# com.flighttrack.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.flighttrack.daily</string>
  <key>ProgramArguments</key>
  <array>
    <string>FLIGHT_TRACK_HOME/.venv/bin/python</string>
    <string>FLIGHT_TRACK_HOME/tracker.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>FLIGHT_TRACK_HOME</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HEADLESS</key>
    <string>1</string>
  </dict>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>FLIGHT_TRACK_LOG_DIR/flight-track-launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>FLIGHT_TRACK_LOG_DIR/flight-track-launchd.err.log</string>
</dict>
</plist>
